# deploy/ecs/agent-runtime.service
[Unit]
Description=Agent Runtime
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/Projects/web/backend
EnvironmentFile=/home/deploy/Projects/web/backend/.env
ExecStart=/home/deploy/miniconda3/envs/backend/bin/python -m gunicorn app.agent_runtime.main:app \
  -k uvicorn.workers.UvicornWorker \
  --bind 127.0.0.1:8100 \
  --workers 1 \
  --timeout 600 \
  --access-logfile - \
  --error-logfile -
Restart=always
RestartSec=3
KillSignal=SIGQUIT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
